# dirx-organizer/dirx.service
[Unit]
Description=DirX Organizer - Automatic File Sorting
After=default.target

[Service]
ExecStart=%h/.local/bin/watcher.sh
WorkingDirectory=%h/.local/bin
Restart=always
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
